# dot_config/autostart/Proton Mail Bridge.desktop
[Desktop Entry]
Type=Application
Name=Proton Mail Bridge
Exec="ch.protonmail.protonmail-bridge" "--noninteractive"
X-GNOME-Autostart-enabled=true
