# dot_config/autostart/rclone.desktop
[Desktop Entry]
Type=Application
Name=rclone
Exec="rclone" "mount" "pcloud:" "/home/radim/pCloudDrive"
X-GNOME-Autostart-enabled=true
